# sync/crontab
# run calendar_sync.sh every minute
* * * * * /usr/local/bin/calendar_sync.sh >> /var/log/cron.log 2>&1

